# .docker/supervisor/nginx.ini
[program:nginx]
command = nginx
startsecs = 2

; Redirect output to syslog
stdout_syslog = true
stderr_syslog = true

; FPM > NGINX
priority = 500
